# templates/etc/supervisor/conf.d/elasticsearch.conf.j2
[program:elasticsearch]
command=/usr/share/elasticsearch/bin/elasticsearch
user = elasticsearch
priority = 999
startsecs = 5
redirect_stderr = True
autostart = True
